# casdoor.service
[Unit]
Description=Casdoor Authentication Server
After=network.target

[Service]
Type=simple
User=jasper
WorkingDirectory=/mnt/usb/projects/casdoor
ExecStart=/mnt/usb/projects/casdoor/casdoor
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
