# logto-auth.service
[Unit]
Description=Logto Authentication Server
After=network.target

[Service]
Type=simple
User=jasper
WorkingDirectory=/mnt/usb/projects/webapp/auth-server
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment="NODE_ENV=production"
Environment="PORT=4000"

[Install]
WantedBy=multi-user.target
