# webapp.service
[Unit]
Description=TypeScript Web App
After=network.target

[Service]
Type=simple
User=jasper
WorkingDirectory=/mnt/usb/projects/webapp
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=10
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target

# To install:
# sudo cp webapp.service /etc/systemd/system/
# sudo systemctl enable webapp
# sudo systemctl start webapp